# crontab
*/5 * * * * /opt/venv/bin/python /Scraper.py >> /logfile.log 2>&1 
